# packaging/debian/postinst
#!/bin/sh
set -eu

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload >/dev/null 2>&1 || true
fi

echo ""
echo "AI Agent Bridge installed."
echo ""
echo "To start the bridge automatically when you log in, run:"
echo "  systemctl --user enable --now bridge"
echo ""
echo "To start it immediately for this session:"
echo "  systemctl --user start bridge"
echo ""
echo "Check status with:"
echo "  bridgectl server status"
echo ""
