# packaging/debian/prerm
#!/bin/sh
set -eu

# The bridge runs as a user service; per-user enablement is the user's
# responsibility. We do not touch user-unit state on package removal.
# Users can stop and disable the service themselves:
#   systemctl --user stop bridge
#   systemctl --user disable bridge
