# packaging/bridge.user.service
[Unit]
Description=AI Agent Bridge (user session)
Documentation=https://github.com/markcallen/ai-agent-bridge
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/bridgectl server start
Restart=on-failure
RestartSec=5s
Environment=HOME=%h

[Install]
WantedBy=default.target
